# conf/readflow.cron
0 * * * * __INSTALL_DIR__/readflow sync -c __INSTALL_DIR__/readflow.yaml 2>> /var/log/__APP__/__APP__.log
