# conf/readflow.cron
0 * * * * __INSTALL_DIR__/readflow sync -c __INSTALL_DIR__/readflow.yaml sync 2>> /var/log/"$app"/"$app".log
